# Docker/systemd_playground/systemd/foo2.service
[Unit]
Description=Test service foo 2
StopWhenUnneeded=true

[Service]
ExecStart=/bin/bash /usr/local/bin/foo.sh 2
Type=simple


