# Docker/systemd_playground/systemd/groupA.target
[Unit]
Description=Test target master
Wants=foo1.service foo2.service
Conflicts=groupB.target

[Install]
WantedBy=multi-user.target
