# Docker/systemd_playground/systemd/groupB.target
[Unit]
Description=Test target slave
Wants=foo2.service foo3.service
Conflicts=groupA.target
