# Docker/systemd_playground/systemd/groupB.target
[Unit]
Description=GroupB target
Wants=foo2.service foo3.service
Conflicts=groupA.target
