# Docker/systemd_playground/systemd/master.target
[Unit]
Description=Test target master
Wants=foo1.service foo2.service
Conflicts=slave.target

[Install]
WantedBy=multi-user.target
